# docker/etc/cron.d/crontest
* * * * * root (date) >> /var/www/.crontest.log
